# hostingservices/memcache_services/templates/memcache_services/supervisor.conf
[program:{{service_name}}]
command=memcached -m {{service.memory_size}} -p {{port}}
user=nobody
autostart=true
autorestart=true
redirect_stderr=True
